# services/apt-update.service
[Unit]
Description=APT Updater
Wants=apt.update.timer

[Service]
Type=simple
ExecStart=/usr/bin/apt update
