# services/apt-update.timer
[Unit]
Description=Run APT update every 3 hours

[Timer]
OnBootSec=3h
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
